# systemd/no-bla-me-dev.service
[Unit]
Description=No-Bla-Me Web Application - Development
Documentation=https://github.com/pbertain/no-bla-me-web
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/dev/no-bla-me
ExecStart=/usr/bin/python3 -m http.server 19080 --bind 0.0.0.0
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=no-bla-me-dev

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dev/no-bla-me
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target